# backend/crontab
# Run the ship_data cleanup task every day at 3 AM (deletes data older than 3 months)
0 3 * * * cd /app && python manage.py cleanup_ship_data >> /var/log/cron.log 2>&1
# Empty line at end of file is required for cron 